# deploy/systemd/starweft-principal.service
[Unit]
Description=Starweft Principal
After=network-online.target starweft-relay.service starweft-owner.service
Wants=network-online.target

[Service]
Type=simple
User=starweft
Group=starweft
EnvironmentFile=-/etc/starweft/principal.env
WorkingDirectory=/var/lib/starweft/principal
ExecStart=/usr/local/bin/starweft run --data-dir /var/lib/starweft/principal --foreground
Restart=always
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/starweft/principal

[Install]
WantedBy=multi-user.target
